# server/whisper_service/whisper.service
[Unit]
Description=LazusAI Whisper transcription service
After=network.target

[Service]
Type=simple
WorkingDirectory=__SERVICE_DIR__
# Bind to localhost only — internal endpoint reached by n8n on the same host.
# WHISPER_MODEL, LAZUSAI_WHISPER_KEY, BLUEBUBBLES_URL, BLUEBUBBLES_PASSWORD
# can be supplied via /etc/lazusai/whisper.env
EnvironmentFile=-/etc/lazusai/whisper.env
Environment=WHISPER_MODEL=base
ExecStart=__SERVICE_DIR__/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 8002
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
